# pi/rc.local
sudo bash /home/paytag/Desktop/paytag/run.sh
sudo bash /home/paytag/Desktop/kiosk/run.sh

exit 0